# ansible/k3s/roles/cleanup/files/k3s.service.d/containerd-mount.conf
[Unit]
# Ensure k3s doesn't start until containerd cache is mounted
RequiresMountsFor=/var/lib/rancher/k3s/agent/containerd
After=var-lib-rancher-k3s-agent-containerd.mount
