# ansible/k3s/roles/cleanup/files/k3s.service.d/dependencies.conf
[Unit]
Requires=k3s-config-init.service admission-controller.service k3s-cluster-init.service
After=k3s-config-init.service admission-controller.service
Wants=k3s-cluster-init.service

[Service]
# Check that config exists (it should after k3s-config-init runs)
# This will fail if k3s-config-init failed to create the config
ExecStartPre=/bin/bash -c 'test -f /etc/rancher/k3s/config.yaml || (echo "ERROR: k3s config not found" && exit 1)'