# ansible/k3s/roles/cleanup/files/k3s.service.d/k3s-dependencies.conf
[Unit]
# Ensure k3s-init runs before k3s
After=k3s-init.service
Requires=k3s-init.service

# Also wait for network to be fully ready
After=network-online.target
Wants=network-online.target

[Service]
# Check that config exists (it should after k3s-init runs)
# This will fail if k3s-init failed to create the config
ExecStartPre=/bin/bash -c 'test -f /etc/rancher/k3s/config.yaml || (echo "ERROR: k3s config not found" && exit 1)'