# ansible/k3s/roles/security/files/module-monitor/module-monitor.service
[Unit]
Description=Module Security Monitor
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/module-monitor.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target