# ansible/k3s/roles/cache-volume/templates/verify-storage.service.j2
[Unit]
Description=Verify Storage Volume Encryption
Documentation=https://github.com/yourusername/sek8s
After=cache-storage.mount
Before=setup-storage-bind-mounts.service containerd.service k3s.service k3s-agent.service
Requires=cache-storage.mount
DefaultDependencies=no
# If this service fails, don't continue boot
OnFailure=poweroff.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Don't timeout - verification is critical
TimeoutStartSec=infinity
# Output to both journal and console for debugging
StandardOutput=journal+console
StandardError=journal+console

# Pass debug mode flag to script
Environment="DEBUG_MODE={{ 'true' if debug_build | default(false) else 'false' }}"

# Verify storage encryption before any data copying
ExecStart=/usr/local/bin/verify-storage.sh

[Install]
WantedBy=local-fs.target
