# ansible/k3s/roles/security/files/module-attestation/module-attestation.timer
[Unit]
Description=Run Module Attestation Every 5 Minutes
Requires=module-attestation.service

[Timer]
# Run every 5 minutes
OnBootSec=2min
OnUnitActiveSec=5min
# Randomize by up to 30 seconds to prevent thundering herd
RandomizedDelaySec=30
# Persist across reboots
Persistent=true

[Install]
WantedBy=timers.target