# ansible/k3s/roles/security/files/module-attestation/module-attestation.service
[Unit]
Description=Module Security Attestation Collection
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/collect-module-state.sh
StandardOutput=journal
StandardError=journal
# Retry on failure
Restart=on-failure
RestartSec=60s
# Maximum 3 attempts
StartLimitBurst=3
StartLimitInterval=300

[Install]
WantedBy=multi-user.target