# ansible/k3s/roles/system-status/files/system-status.service
[Unit]
Description=sek8s System Status Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=status
Group=status
WorkingDirectory=/opt/sek8s
EnvironmentFile=/etc/system-status/system-status.env
ExecStart=/opt/sek8s/venv/bin/python -m sek8s.services.system_status
Restart=always
RestartSec=5
RuntimeDirectory=system-status
RuntimeDirectoryMode=0750
StandardOutput=journal
StandardError=journal
SyslogIdentifier=system-status

[Install]
WantedBy=multi-user.target
