# ansible/k3s/roles/gpu/files/infiniband-mask.service
[Unit]
Description=Mask InfiniBand services when no Mellanox/ConnectX devices are present
Documentation=See config-volume-validator.service for device-dependent service pattern
# Run after udev settles so passthrough PCI devices are visible (align with config-volume-validator)
After=systemd-udev-settle.service
After=local-fs.target
Before=openibd.service opensmd.service rdma-core.service rdma.service ibacm.service
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/infiniband-mask.sh
TimeoutStartSec=30
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
