# ansible/k3s/roles/admission-controller/files/admission-controller.service
[Unit]
Description=Kubernetes Admission Controller
After=network.target opa.service
Wants=opa.service
Before=k3s.service

[Service]
Type=simple
User=admission
Group=admission
WorkingDirectory=/opt/admission-controller

# Use virtual environment
Environment="PATH=/opt/admission-controller/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/opt/admission-controller/venv"
EnvironmentFile=/etc/admission-controller/admission.env

# Run the admission controller
ExecStart=/opt/admission-controller/venv/bin/python -m sek8s.admission_controller

Restart=always
RestartSec=5
StartLimitInterval=300
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=admission-controller

[Install]
WantedBy=multi-user.target