# ansible/k3s/roles/admission-controller/files/opa.service
[Unit]
Description=Open Policy Agent
After=network.target
Before=admission-controller.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/opa run --server \
    --config-file=/etc/opa/opa.yaml \
    --addr=localhost:8181 \
    --diagnostic-addr=0.0.0.0:8282 \
    --log-level=info \
    /etc/opa/policies
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/opa
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

[Install]
WantedBy=multi-user.target