# ansible/k3s/roles/cache-volume/files/setup-chutes-agent-symlink.service
[Unit]
Description=Setup symlink for Chutes agent state to cache volume
Documentation=man:systemd.service(5)
Requires=var-snap.mount
After=var-snap.mount
Before=k3s.service
Before=k3s-agent.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/setup-chutes-agent-symlink.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
